# etc/supervisor/conf.d/copy-authorized-keys.conf
[program:copy-authorized_keys]
command=/opt/copy-authorized_keys.sh
numprocs=1
autostart=true
autorestart=false
